# schedules/com.eathy.xhs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eathy.xhs</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/zeke/agent dome/eathy-xhs/.venv/bin/eathy</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/zeke/agent dome/eathy-xhs</string>

    <key>StartCalendarInterval</key>
    <array>
        <!-- 08:00 -->
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 12:00 -->
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 18:00 -->
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zeke/agent dome/eathy-xhs/data/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zeke/agent dome/eathy-xhs/data/logs/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
